# arp2bgp/arp-2-bgp.cron
# Runs the arp-2-bgp script
*/5 * * * *  /usr/lib/python2.7/site-packages/arp2bgp/arp-2-bgp.py /etc/arp-2-bgp.conf > /dev/null 2>&1

