# arp2bgp/arp-2-bgp.cron
# Runs the arp-2-bgp script
*/5 * * * *  /usr/local/bin/arp-2-bgp.py /etc/arp-2-bgp.conf >> /tmp/arp-2-bgp.log

